# etc/cron.daily/backup-local
#!/bin/bash

[ "$1" != "--really" ] || exec nice ionice "$0" --really

HOSTNAME=$( hostname -s )
BACKUP="/slack/backups/$HOSTNAME"
USER=$( getent passwd 1000 | cut -d: -f1 )

if [ -z "$USER" ]; then
  echo "Couldn't find default user"
  exit 1
fi

function user {
  sudo -u "$USER" sh -c "$1"
}

umask 0027
mkdir -p "$BACKUP" || exit 1
cd "$BACKUP" || exit 1
chown "$USER:$USER" "$BACKUP"

user "dpkg --get-selections > packages.txt"

if [ -d "/home/$USER/.config/dconf" ]; then
  user "dconf dump / > dconf.ini"
fi

for dir in movies series; do
  [ -d "/slack/$dir" ] && user "find /slack/$dir/ -type f | sort > $dir.txt"
done
