# etc/cron.daily/clear-tempfiles
#!/bin/bash

for path in {/home/*,/root}/{.cache/nvim,.local/{share,state}/nvim/{swap,undo}}; do
  if [ -d "$path" ]; then
    find "$path" -type f -mtime +90 -delete
  fi
done
