# etc/cron.daily/clear-tempfiles
#!/bin/bash

find /tmp \
  -maxdepth 1 -mtime +1 \( \
  -name 'dropbox-*' -o \
  -name 'gitaly-*' -o \
  -name 'karma-*' -o \
  -name 'ruby-*' -o \
  -name 'xvfb-run.*' -o \
  -name '.com.google.*' -o \
  -name '.org.chromium.*' \
  \) -exec rm -rf "{}" \;

find /home/*/{.cache/nvim,.local/share/nvim/{swap,undo}} \
  -type f -mtime +90 -delete 2>/dev/null
