# etc/profile.d/environment.sh
export EDITOR="sensible-vim"
export GIT_EDITOR="vim"
export PAGER="less"

export LESS="-iRM"
export GREP_OPTIONS="-i --color=auto --exclude=.svn --exclude=.git --exclude=.*.swp"
export ACK_OPTIONS="--smart-case --ignore-dir=tmp/cache --ignore-dir=tmp/files"
export IRB_HISTORY_SIZE=3000
export RI="--format ansi --no-pager"
