# etc/profile.d/environment-fnord.sh
export VDPAU_NVIDIA_SYNC_DISPLAY_DEVICE="DFP-1"
#export __GL_SYNC_TO_VBLANK="1"
