# etc/profile.d/environment.sh
export PATH=/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/sbin:/usr/local/bin:/usr/games:/var/lib/gems/1.8/bin:~/bin

export LANG="de_CH.UTF-8"
export LC_MESSAGES="en_US.UTF-8"
export LC_NUMERIC="en_US.UTF-8"
export LC_CTYPE="en_US.UTF-8"

export EDITOR="vim"
export PAGER="less"

export LESS="-iRM"
export GREP_OPTIONS="-i --color=auto --exclude=.svn --exclude=.git --exclude=.*.swp"
export ACK_OPTIONS="-i"
export IRB_HISTORY_SIZE=3000
export RI="--format ansi --no-pager"
