# etc/profile.d/nodejs.sh
# Check for interactive bash
[ -n "$BASH_INTERACTIVE" ] || return

alias sunpm='sudo npm -g'
