# etc/profile.d/z99-fnord.sh
export CDPATH=.:~:/etc/dotfiles:~/politnetz

export LANG="de_CH.UTF-8"
export LC_MESSAGES="en_US.UTF-8"
export LC_NUMERIC="en_US.UTF-8"
export LC_CTYPE="en_US.UTF-8"

export VDPAU_NVIDIA_SYNC_DISPLAY_DEVICE="DFP-1"
#export __GL_SYNC_TO_VBLANK="1"
