# etc/udev/rules.d/99-steamdeck-hdmi.rules
ACTION=="add", SUBSYSTEM=="cec", RUN+="/slack/dotfiles/bin/switch-hdmi"
