# etc/udev/rules.d/99-steamdeck-hdmi.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/usr/bin/systemd-run --no-block /slack/dotfiles/bin/switch-hdmi"
