# config/autostart/setup-peripherals.desktop
[Desktop Entry]
Name=setup-peripherals
Exec=/slack/dotfiles/bin/setup-peripherals
Type=Application
