# config/autostart/setup-peripherals.desktop
[Desktop Entry]
Name=setup-peripherals
Exec=setup-peripherals
Type=Application
